# docker/php/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm

[program:laravel-scheduler]
command=php /var/www/artisan schedule:work
autostart=true
autorestart=true