# docker/php/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null ; ログファイルを無効化
pidfile=/tmp/supervisord.pid ; PIDファイルを/tmpに変更

[program:php-fpm]
command=/usr/local/bin/php
autostart=true
autorestart=true
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm.err

[program:cron]
command=cron -f
autostart=true
autorestart=true
stdout_logfile=/var/log/cron.log
stderr_logfile=/var/log/cron.err